# emoncms/rootfs/etc/services.d/php-fpm/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: Emoncms
# Runs the PHP-FPM daemon
# ==============================================================================
export MYSQL_HOST
export MYSQL_NAME
export MYSQL_USERNAME
export MYSQL_PASSWORD
export MYSQL_PORT

if bashio::config.has_value 'remote_mysql_host';then
    MYSQL_HOST=$(bashio::config "remote_mysql_host")
    MYSQL_NAME=$(bashio::config "remote_mysql_database")
    MYSQL_USERNAME=$(bashio::config "remote_mysql_username")
    MYSQL_PASSWORD=$(bashio::config "remote_mysql_password")
    MYSQL_PORT=$(bashio::config "remote_mysql_port")
else
    MYSQL_HOST=$(bashio::services "mysql" "host")
    MYSQL_NAME=emoncms
    MYSQL_USERNAME=$(bashio::services "mysql" "username")
    MYSQL_PASSWORD=$(bashio::services "mysql" "password")
    MYSQL_PORT=$(bashio::services "mysql" "port")
fi

bashio::log.info "Starting PHP-FPM server..."
exec php-fpm7 --nodaemonize
